# nico.service
[Unit]
Description=NICO AI Assistant
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=yatin
WorkingDirectory=/home/yatin/nico
ExecStart=/home/yatin/nico/.venv/bin/python /home/yatin/nico/main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
