# nico.service
[Unit]
Description=NICO AI Assistant — Smart Speaker
After=network-online.target bluetooth.target sound.target
Wants=network-online.target bluetooth.target

[Service]
Type=simple
User=yatin
WorkingDirectory=/home/pi/nico
ExecStart=/home/pi/nico/.venv/bin/python /home/pi/nico/main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

# Bluetooth audio — wait for PulseAudio
ExecStartPre=/usr/bin/sleep 5

[Install]
WantedBy=multi-user.target
