# packaging/debian/nico.service
[Unit]
Description=NICO AI Assistant
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=nico
Group=nico
WorkingDirectory=/opt/nico
Environment=NICO_ENV_FILE=/etc/nico/.env
ExecStart=/opt/nico/.venv/bin/python /opt/nico/main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardInput=tty
TTYPath=/dev/tty1

[Install]
WantedBy=multi-user.target
